# trinity-brain.service
[Unit]
Description=Trinity AI Brain - vLLM Mistral Small 4 Inference Server (P-EAGLE)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=joshua
WorkingDirectory=/home/joshua/Workflow/desktop_trinity/trinity-genesis

# Environment for AMD Strix Halo / RDNA 3.5
Environment="HSA_OVERRIDE_GFX_VERSION=11.0.0"
Environment="ROCBLAS_USE_HIPBLASLT=1"
Environment="HIP_VISIBLE_DEVICES=0"
Environment="AMD_LOG_LEVEL=0"
# Bind vLLM to 8080 so the Rust backend connects natively without codebase refactoring
Environment="VLLM_PORT=8080"  

# ═══════════════════════════════════════════════════════════════════════════════
# vLLM MS4 + P-EAGLE Launcher
# ═══════════════════════════════════════════════════════════════════════════════

# Start the vLLM Brain
ExecStart=/home/joshua/Workflow/desktop_trinity/trinity-genesis/scripts/launch/start_vllm_ms4.sh

# Restart policy for stability
Restart=always
RestartSec=10

# Memory limits (Strix Halo UMA ceiling)
MemoryMax=110G

# Graceful shutdown timeout (allows vLLM to flush queues)
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
