# trinity-brain.service
[Unit]
Description=Trinity AI Brain - Embedded llama.cpp Inference (Vulkan GPU)
After=network.target

[Service]
Type=simple
User=joshua
WorkingDirectory=/home/joshua/Workflow/desktop_trinity/trinity-genesis

# Inference architecture: Embedded llama-cpp-2 → Vulkan GPU
# No separate server process — Trinity loads the GGUF directly
Environment="GGUF_MODEL_PATH=/home/joshua/trinity-models/gguf/Mistral-Small-4-119B-2603-Q4_K_M-00001-of-00002.gguf"

# Strix Halo Vulkan fix (>64GB models)
Environment="RADV_PERFMODE=nogttspill"

# Start Trinity (which loads the model via embedded inference)
ExecStart=/home/joshua/Workflow/desktop_trinity/trinity-genesis/target/release/trinity

# Restart policy for stability
Restart=on-failure
RestartSec=15

# Memory limits (Strix Halo UMA ceiling)
MemoryMax=110G

# Graceful shutdown timeout
TimeoutStopSec=30

[Install]
WantedBy=default.target
