# scripts/systemd/trinity-llama.service
[Unit]
Description=Trinity Llama Server (Mistral 119B)
# Start early, before anything else in the Trinity stack
After=network.target
Wants=network.target
# Prevent anything from starting until this is fully up
Before=trinity-os.service

[Service]
Type=simple
User=joshua
WorkingDirectory=/home/joshua/Workflow/desktop_trinity/trinity-genesis

ExecStart=/home/joshua/Workflow/desktop_trinity/trinity-genesis/llama.cpp/build-vulkan/bin/llama-server \
    --model /home/joshua/trinity-models/gguf/Mistral-Small-4-119B-2603-Q4_K_M-00001-of-00002.gguf \
    --port 8080 --host 127.0.0.1 \
    --ctx-size 262144 --n-gpu-layers 99 --flash-attn on --jinja --parallel 2

# ══ Memory Pinning — "don't touch this process" ══
# Tell the OOM killer to NEVER kill llama (it's our most important process)
OOMScoreAdjust=-900
# Give it higher CPU priority during model loading
Nice=-5
# Lock memory pages so they can't be swapped to disk
LimitMEMLOCK=infinity

Restart=on-failure
RestartSec=10
# Model takes ~60-90s to load into VRAM — don't time out early
TimeoutStartSec=180
Environment=HOME=/home/joshua

[Install]
WantedBy=multi-user.target
