# scripts/systemd/trinity-voxtral.service
[Unit]
Description=Trinity Voxtral — 4B TTS Engine
Documentation=https://huggingface.co/mistralai/Voxtral-4B-TTS-2603
After=trinity-vllm.service
Wants=trinity-vllm.service

[Service]
Type=exec
User=joshua
Group=joshua
WorkingDirectory=/home/joshua/Workflow/desktop_trinity/trinity-genesis

# AMD Strix Halo ROCm
Environment=HSA_OVERRIDE_GFX_VERSION=11.0.0
Environment=ROCBLAS_USE_HIPBLASLT=1
Environment=HIP_VISIBLE_DEVICES=0
Environment=AMD_LOG_LEVEL=0

ExecStart=/home/joshua/Workflow/desktop_trinity/trinity-genesis/scripts/launch/start_voxtral.sh

# Graceful shutdown
KillSignal=SIGINT
TimeoutStopSec=15

# Auto-restart — TTS is non-critical, delay longer
Restart=on-failure
RestartSec=30

# Lower memory — 4B model needs ~6 GB
MemoryMax=16G
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
