# scripts/systemd/trinity-vllm.service
[Unit]
Description=Trinity vLLM — Mistral Small 4 (119B MoE) Inference Server
Documentation=https://github.com/vllm-project/vllm
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=joshua
Group=joshua
WorkingDirectory=/home/joshua/Workflow/desktop_trinity/trinity-genesis

# AMD Strix Halo ROCm
Environment=HSA_OVERRIDE_GFX_VERSION=11.0.0
Environment=ROCBLAS_USE_HIPBLASLT=1
Environment=HIP_VISIBLE_DEVICES=0
Environment=AMD_LOG_LEVEL=0

ExecStart=/home/joshua/Workflow/desktop_trinity/trinity-genesis/scripts/launch/start_vllm_ms4.sh

# Graceful shutdown — vLLM flushes pending requests
KillSignal=SIGINT
TimeoutStopSec=30

# Auto-restart on crash
Restart=on-failure
RestartSec=10

# Resource limits — Strix Halo 128 GB UMA
MemoryMax=100G
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
